# livecd/trunk/root_overlay/etc/local.d/99-remount-cd-rw.start
#!/bin/sh

# workaround read-only mount
#this should be fixed in pm-utils
mount -o remount,rw /mnt/cdrom
